# systemd/hud.service
[Unit]
Description=Car HUD System
After=graphical.target network.target
Wants=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/car-hud-pi
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/pi/.Xauthority
ExecStartPre=/bin/sleep 5
ExecStart=/home/pi/car-hud-pi/scripts/start.sh
Restart=always
RestartSec=5

[Install]
WantedBy=graphical.target
