# systemd/hud.service
[Unit]
Description=Car HUD System
After=graphical.target network.target
Wants=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/car-hud-pi
Environment=XAUTHORITY=/home/pi/.Xauthority
# No HDMI/kiosk: keep FastAPI + WebSocket up without Chromium (see scripts/start.sh).
# Environment=HUD_SKIP_CHROMIUM=1
ExecStartPre=/bin/sleep 5
ExecStart=/home/pi/car-hud-pi/scripts/start.sh
Restart=always
RestartSec=5

[Install]
WantedBy=graphical.target
